# hyper-bot.service
[Unit]
Description=Hyperliquid Trading Bot
After=network.target

[Service]
Type=simple
User=%USER%
Group=%USER%
WorkingDirectory=/home/%USER%/hyper-bot
Environment="PATH=/home/%USER%/hyper-bot/venv/bin"
ExecStart=/home/%USER%/hyper-bot/venv/bin/python main.py
Restart=always
RestartSec=10

# Logs (créer le dossier avant de démarrer)
StandardOutput=append:/home/%USER%/hyper-bot/logs/bot.log
StandardError=append:/home/%USER%/hyper-bot/logs/bot_error.log

# S'assurer que le dossier logs existe et a les bonnes permissions
ExecStartPre=/bin/mkdir -p /home/%USER%/hyper-bot/logs
ExecStartPre=/bin/chown -R %USER%:%USER% /home/%USER%/hyper-bot/logs
ExecStartPre=/bin/chmod 755 /home/%USER%/hyper-bot/logs

# Sécurité (optionnel)
# NoNewPrivileges=true
# PrivateTmp=true

[Install]
WantedBy=multi-user.target

